# service_templates/records.timer
[Unit]
Description="Sync OpenStack with SUPR"

[Timer]
OnCalendar=*,37
Unit=records.service

[Install]
WantedBy=multi-user.target
