# service_templates/supr_sync.timer
[Unit]
Description="Sync OpenStack with SUPR"

[Timer]
OnCalendar=*:0,15,30,45
Unit=supr-sync.service

[Install]
WantedBy=multi-user.target
